# launchagents/com.vespera.api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.vespera.api</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>REPLACE_WITH_YOUR_PATH/vespera/api.py</string>
    </array>
    <!-- Note: api.py auto-detects gunicorn and uses it if installed.
         After running setup.sh, gunicorn will be available in your venv.
         No plist changes needed — the script handles it. -->
    <key>WorkingDirectory</key>
    <string>REPLACE_WITH_YOUR_PATH/vespera</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/vespera-api.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/vespera-api.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
</dict>
</plist>
